# docker/php/crontab
*/5 * * * * /usr/local/bin/php /var/www/bin/console mevo:fetch >> /var/log/cron.log 2>&1
